# .config/systemd/user/jack.service
[Unit]
Description=JACK
After=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/jackd -R -d alsa 

[Install]
WantedBy=default.target
