# .config/systemd/user/pulseaudio.service
[Unit]
Description=PulseAudio
After=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/pulseaudio

[Install]
WantedBy=default.target
